# l/bash/setup/cntshjtg.conf
[program:control]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:control_1]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT