# test/setup/all.conf
[program:entrypoint_105]
command=/scripts/%(program_name)s.sh
log_stdout=false
log_stderr=false
redirect_stderr=true
autostart=true
autorestart=false
stopasgroup=true
stopsignal=INT

[program:entrypoint_106]
command=/scripts/%(program_name)s.sh
log_stdout=false
log_stderr=false
redirect_stderr=true
autostart=true
autorestart=false
stopasgroup=true
stopsignal=INT

[program:agent]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
startretries=5
stopasgroup=true
stopsignal=INT

[program:keeper]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:fms_adaptor]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:vehicle]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
stopasgroup=true
startretries=5
stopsignal=INT
priority=10

[program:control]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:identification]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:gnss_driver]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
startretries=5
stopasgroup=true
stopsignal=INT
priority=100

[program:gnss_processor]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
startretries=5
stopasgroup=true
stopsignal=INT

[program:wheel_odom]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:localization_adaptor]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:localization_checker]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
startretries=5
stopasgroup=true
stopsignal=INT
priority=999

[program:localization_logger]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
stopasgroup=true
startretries=5
stopsignal=INT
priority=999

[program:fusion_localizer2]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
autorestart=true
startretries=5
stopasgroup=true
stopsignal=INT
priority=999

[program:video]
command=/scripts/%(program_name)s.sh
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
startretries=5
stopasgroup=true
stopsignal=INT

[program:video_1]
command=/scripts/video-ck.sh 0
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
startretries=5
stopasgroup=true
stopsignal=INT

[program:http_bridge]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
startretries=5
stopasgroup=true
stopsignal=INT

[program:aeb]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:local_plan]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:vdr]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
autorestart=false
autostart=true
startretries=5
stopasgroup=true
stopsignal=INT

[program:test_rviz]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
startretries=5
stopasgroup=true
stopsignal=INT

[program:test_vehicle_sim_plan]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:lidar_config_check]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
startretries=5
autostart=false
stopasgroup=true
stopsignal=INT

[program:qp3_lidarpreprocess]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
autorestart=true
stopasgroup=true
startretries=5
stopsignal=INT
priority=999

[program:qp3_lidarpreprocess_1]
command=/scripts/qp3_perception.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
autostart=false
autorestart=true
stopasgroup=true
startretries=5
stopsignal=INT
priority=999

[program:qp3_lidarpreprocess_2]
command=/scripts/qp3_transitional_perception.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:qp3_perception]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:landmark_localizer_ros2]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=20
redirect_stderr=true
autostart=false
autorestart=true
startretries=5
stopasgroup=true
stopsignal=INT
priority=999

[program:motion_estimator]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:alignment_planner]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:lidar_cps_alignment]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:remote_control_adaptor]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:remote_control]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT

[program:lidar_control]
command=/scripts/%(program_name)s.sh
stdout_logfile=/debug/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=50
redirect_stderr=true
autostart=false
stopasgroup=true
stopsignal=INT
